# roles/arch-i3-desktop/files/pacman-bar.service
[Unit]
Description=Get pacman upgrade status for i3status

[Service]
Type=oneshot
# From dots role
ExecStart=/usr/local/bin/pacman-bar
