# roles/arch-i3-desktop/files/wttr-bar.timer
[Unit]
Description=Timer to get wttr.in weather for i3status

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
