# roles/arch-i3-desktop/files/wttr-bar.service
[Unit]
Description=Get wttr.in status for i3status

[Service]
Type=oneshot
ExecStart=%h/bin/wttrbar
